# etc/udev/rules.d/60-usb-dirty-pages-udev.rules
ACTION=="add", KERNEL=="sd[a-z]", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/usr/bin/usb-dirty-pages-udev %k"
